# lib/lisp/proj2/Paper/cronjob.cron
0,15,30,45 * * * * "~/svns/ourmine/branches/13/lib/lisp/proj2/Paper/make update"
