# lib/lisp/proj2/Paper/cronjob.cron
0,10,20,30,40,50 * * * * "~/svns/ourmine/branches/13/lib/lisp/proj2/Paper/make update"
